# supervisord.conf
[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/workspace/logs/nginx-error.log
stdout_logfile=/workspace/logs/nginx-access.log
priority=100

[program:comfyui]
command=/workspace/ComfyUI/venv/bin/python main.py --listen 0.0.0.0 --port %(ENV_COMFYUI_PORT)s
directory=/workspace/ComfyUI
environment=LD_PRELOAD="$(ldconfig -p | grep -Po 'libtcmalloc.so.\d' | head -n 1)"
autostart=true
autorestart=true
stderr_logfile=/workspace/logs/comfyui-error.log
stdout_logfile=/workspace/logs/comfyui-access.log
priority=200
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface